# packaging/vendor/rak/mipsel_24kc/files/rak-basicstation.init
#!/bin/sh /etc/rc.common

START=99
STOP=10

USE_PROCD=1

PROG=/usr/bin/rak-basicstation

start_service() {
    procd_open_instance
    procd_set_param command "$PROG" -c /etc/rak-basicstation/rak-basicstation.toml
    procd_set_param respawn
    procd_close_instance
}
